# 2LabsToGo-Eco-Backup/backup_sd.service
[Unit]
Description="Runs my demo backup sh script"

[Service]
ExecStart=/etc/systemd/system/backup_sd.sh

