# 2LabsToGo-Eco-Backup/backup_db.service
[Unit]
Description="Runs my demo backup sh script"

[Service]
ExecStart=/etc/systemd/system/backup_db_media.sh

